# deploy/systemd/modelport.service
[Unit]
Description=ModelPort local Claude model gateway
Documentation=https://github.com/tiammomo/ModelPort
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/modelport/modelport.env
Environment=RUST_LOG=model_port=info,tower_http=info
Environment=MODELPORT_STATE_DIR=/var/lib/modelport
Environment=MODELPORT_CONTROL_STORE=/var/lib/modelport/control-plane.json
ExecStart=/usr/local/bin/model-port
WorkingDirectory=/var/lib/modelport
Restart=always
RestartSec=3

DynamicUser=true
StateDirectory=modelport
StateDirectoryMode=0700
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
CapabilityBoundingSet=
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target
